# claude-proxy.service
[Unit]
Description=Claude Code Proxy
After=network.target

[Service]
ExecStart=%h/claude-cli-proxy/venv/bin/python3 %h/claude-cli-proxy/claude-proxy.py
WorkingDirectory=%h/claude-cli-proxy
Restart=always
RestartSec=5
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
